# deployment/systemd/alert_monitor.service
[Unit]
Description=Crypto Monitor Alert Service
After=network.target redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/v8.3_crypto_monitor/redis_server
ExecStart=/usr/bin/python3 /root/v8.3_crypto_monitor/redis_server/alert_monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
