# deployment/systemd/dashboard.service
[Unit]
Description=Crypto Monitor Dashboard
After=network.target crypto-monitor.service
Wants=crypto-monitor.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/crypto-monitor
Environment=PYTHONPATH=/root/crypto-monitor/src
EnvironmentFile=/root/crypto-monitor/.env

ExecStart=/root/crypto-monitor/venv/bin/python -m src.dashboards.unified.app

Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=crypto-dashboard

[Install]
WantedBy=multi-user.target
