# deployment/systemd/fusion-dashboard.service
[Unit]
Description=Fusion Dashboard v9.5
After=network.target redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/fusion_dashboard_v95
EnvironmentFile=/etc/fusion_dashboard.env
ExecStart=/usr/bin/python3 /root/fusion_dashboard_v95/server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
