# deployment/systemd/fusion_engine.service
[Unit]
Description=Crypto Monitor Fusion Engine v3
After=network.target redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/v8.3_crypto_monitor/redis_server
ExecStart=/usr/bin/python3 -u /root/v8.3_crypto_monitor/redis_server/fusion_engine_v3.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
