# deployment/systemd/signal_router.service
[Unit]
Description=Crypto Signal Router
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/v8.3_crypto_monitor/redis_server
ExecStart=/usr/bin/python3 -u /root/v8.3_crypto_monitor/redis_server/signal_router.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
