# deployment/systemd/webhook_pusher.service
[Unit]
Description=Crypto Monitor Webhook Pusher
After=network.target redis.service fusion_engine.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/v8.3_crypto_monitor/redis_server
ExecStart=/usr/bin/python3 /root/v8.3_crypto_monitor/redis_server/webhook_pusher.py
Restart=always
RestartSec=10
StandardOutput=append:/root/v8.3_crypto_monitor/redis_server/webhook.log
StandardError=append:/root/v8.3_crypto_monitor/redis_server/webhook.log

[Install]
WantedBy=multi-user.target
