# deployment/systemd/wechat_pusher.service
[Unit]
Description=Crypto Monitor WeChat Pusher
After=network.target redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/v8.3_crypto_monitor/redis_server
ExecStart=/usr/bin/python3 -u /root/v8.3_crypto_monitor/redis_server/wechat_pusher.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
